# .ebextensions/supervisord/queue-worker.conf
[program:queue-worker]
command=php /var/app/current/laravel/artisan queue:work --env=production --delay=3 --sleep=3 --tries=3
directory=/var/app/current
stdout_logfile=/var/app/support/logs/queue-worker.log
redirect_stderr=true
autostart=true
autorestart=true
